# debian/switchamba.service
[Unit]
Description=Switchamba — automatic EN/RU/UA keyboard layout switching
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/switchamba
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=graphical-session.target
